# crontab
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.ckuZzs/crontab installed on Wed Dec 24 19:06:18 2014)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.4gOoGe/crontab installed on Wed Dec 10 19:26:26 2014)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
# m h  dom mon dow   command
1 1 * * * /bin/bash /root/denyip.sh
1 2 * * 0 /usr/bin/find /root/mnt/ -mtime +30 -name "log.ipv*" -exec rm -f {} \; 
1 3 1 * * echo  > /root/mnt/log.perf && echo > /root/mnt/log.keepalive
1 4 10 * * /bin/bash /root/topology/refreshGeoLite.sh
0 * * * * sleep `cat /root/mnt/second` && /bin/bash /root/mnt/perf.sh >> /root/mnt/log.perf 2>&1
*/5 * * * * /usr/bin/python /root/mnt/keepalive.py >> /root/mnt/log.keepalive 2>&1
@reboot rm -f /root/mnt/hour && rm -f /root/mnt/run && /bin/bash /root/mnt/perf.sh >> /root/mnt/log.perf 2>&1
@reboot /usr/sbin/ntpdate -u s1a.time.edu.cn
